# ops/systemd/helivex-paper.service
[Unit]
Description=Helivex Paper Trading Node (NautilusTrader)
After=network.target postgresql.service
Documentation=https://github.com/soffy88/helivex
StartLimitBurst=3
StartLimitIntervalSec=120

[Service]
Type=simple
WorkingDirectory=/home/soffy/projects/helivex
EnvironmentFile=/home/soffy/projects/helivex/.env
ExecStart=/home/soffy/projects/helivex/venv/bin/python /home/soffy/projects/helivex/paper/run.py
ExecStartPost=/bin/bash -c 'sleep 2 && echo $MAINPID > /tmp/helivex_paper_node.pid'

# SIGTERM triggers on_stop() → closes all positions gracefully
KillSignal=SIGTERM
TimeoutStopSec=30

# Watchdog: the in-process sd_notify actor (paper/sdwatchdog.py) pings every 30s;
# a wedged engine loop stops pinging and systemd restarts after 120s (generous
# margin over the node's connect/boot before the first ping). A watchdog kill is
# a failure → Restart=on-failure restarts; the kill-switch's clean SIGTERM does not.
WatchdogSec=120

# Restart only on crash (exit != 0). Intentional SIGTERM/SIGINT won't trigger restart.
# This lets the portfolio kill switch stop the node without it bouncing back.
Restart=on-failure
RestartSec=10

StandardOutput=append:/home/soffy/helivex_paper.log
StandardError=append:/home/soffy/helivex_paper.log

[Install]
WantedBy=helivex.target
