# ops/systemd/helivex-ingest-funding-binance.service
[Unit]
Description=Helivex ingest — Binance 8h funding history (via helios-proxy SG egress)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/soffy/projects/helivex
EnvironmentFile=/home/soffy/projects/helivex/.env
# No --start: walks full history each run (idempotent upsert, ~8 calls); self-heals gaps.
ExecStart=/home/soffy/projects/helivex/venv/bin/python /home/soffy/projects/helivex/ops/scripts/fetch_binance_funding.py --symbol BTCUSDT
StandardOutput=append:/home/soffy/helivex_ingest_funding_binance.log
StandardError=append:/home/soffy/helivex_ingest_funding_binance.log
