# ops/systemd/helivex-backup.timer
[Unit]
Description=Helivex daily Postgres backup timer

[Timer]
OnCalendar=*-*-* 02:00:00
# Persistent: run a missed backup on next boot (overnight host-off no longer
# silently skips the dump, the failure mode the audit caught on 2026-06-30).
Persistent=true

[Install]
WantedBy=timers.target
